# Hosaka_Field-Terminal/systemd/hosaka-field-terminal-headless.service
[Unit]
Description=Hosaka Field Terminal Appliance (Headless)
After=network-online.target picoclaw-gateway.service
Wants=network-online.target picoclaw-gateway.service

[Service]
Type=simple
WorkingDirectory=/opt/hosaka-field-terminal
Environment=HOSAKA_STATE_PATH=/home/operator/.hosaka/state.json
Environment=PYTHONUNBUFFERED=1
Environment=HOSAKA_BOOT_MODE=headless
Environment=HOME=/home/operator
ExecStart=/opt/hosaka-field-terminal/.venv/bin/python -m hosaka
Restart=always
RestartSec=2
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
